# units/mongo.service
[Unit]
Description=Mongo
After=docker.service

[Service]
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill mongo-rs0-srv0
ExecStartPre=-/usr/bin/docker rm mongo-rs0-srv0
ExecStartPre=/usr/bin/docker pull korczis/mongod
ExecStart=/usr/bin/docker run -p 127.0.0.1:27017:27017 --name mongo-rs0-srv0 korczis/mongod --replSet rs1 --noprealloc --smallfiles
ExecStop=/usr/bin/docker stop mongo-rs0-srv0